# systemd/trialinfo-secure.socket
[Socket]
ListenStream=443
Service=trialinfo.service

[Install]
WantedBy=sockets.target
